# systemd/computer-use.service.example
# Install into ~/.config/systemd/user/computer-use.service
# Adjust the ExecStart/ExecStop paths if your clone is not at ~/hermes-computer-use.
# Then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now computer-use.service
# Run `sudo loginctl enable-linger $USER` once so it survives logout.

[Unit]
Description=hermes-computer-use virtual display stack (Xvfb + fluxbox + x11vnc + Chrome)
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/bash %h/hermes-computer-use/scripts/display.sh start
ExecStop=/usr/bin/bash %h/hermes-computer-use/scripts/display.sh stop
Environment=CU_DISPLAY=99
Environment=CU_WIDTH=1440
Environment=CU_HEIGHT=900
Environment=CU_VNC_PORT=5900
Environment=CU_STATE_DIR=%t/hermes-computer-use
Environment=CU_PROFILE_DIR=%h/.local/share/hermes-computer-use/chrome-profile

[Install]
WantedBy=default.target
