# systemd/novnc.service
[Unit]
Description=noVNC HTTP proxy (browser-based VNC for computer-use)
After=computer-use.service
Requires=computer-use.service

[Service]
Type=simple
ExecStart=/home/geunuk/.hermes/hermes-agent/venv/bin/websockify \
    --web=/home/geunuk/.hermes/integrations/computer-use/novnc \
    6080 localhost:5900
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
